# roles/restic_profile/templates/etc/systemd/system/restic-profile.timer.j2
{# jinja2: lstrip_blocks: True, trim_blocks: True #}
{% set _profile_name = item.profile_name %}
{% set _type = item.type %}
{% set _unit_name = item.unit_name %}
{% set _sub_config = item.sub_config %}
# /etc/systemd/system/{{ _unit_name }}.timer
# Generated by Ansible — do not edit manually

[Unit]
Description=Restic {{ _type }} timer — {{ _profile_name }}

[Timer]
# Run on this schedule; Persistent=true catches up a missed run after reboot/sleep
OnCalendar={{ _sub_config.on_calendar | default('hourly') }}
Persistent=true
# Spread load: fire at a random offset within this window
RandomizedDelaySec={{ _sub_config.randomized_delay_sec | default('10min') }}

[Install]
WantedBy=timers.target
