# roles/restic_profile/templates/etc/systemd/system/restic-profile.service.j2
{# jinja2: lstrip_blocks: True, trim_blocks: True #}
{% set _profile_name = _profile_item.key %}
{% set _profile = _profile_item.value %}
{% set _is_backup = (_profile.sources | default([]) | length) > 0 %}
{% set _cmd = "backup" if _is_backup else "forget" %}
{% set _unit_prefix = restic_profile_backup_unit_prefix if _is_backup else restic_profile_forget_unit_prefix %}
# {{ _unit_prefix }}{{ _profile_name }}.service
# Generated by Ansible — do not edit manually

[Unit]
Description=Restic {{ _cmd }} — {{ _profile_name }}
# Soft dependency: start after network is up, but don't require it
# (local-only repositories work without a network connection)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
{% if _profile.system_user | default('root') != 'root' %}
User={{ _profile.system_user }}
{% endif %}
ExecStart={{ restic_profile_bin }} {{ _cmd }} {{ _profile_name }} \
    --config {{ restic_profile_config_file }}
# Let systemd capture stdout/stderr in the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ _unit_prefix }}{{ _profile_name }}

[Install]
WantedBy=multi-user.target
