# roles/restic_profile/templates/etc/systemd/system/restic-profile.service.j2
{# jinja2: lstrip_blocks: True, trim_blocks: True #}
{% set _profile_name = _plan_unit_item.profile_name %}
{% set _profile = _plan_unit_item.profile %}
{% set _unit_name = _plan_unit_item.unit_name %}
{% set _cpu_quota = _profile.cpu_quota if _profile.cpu_quota is defined else restic_profile_systemd_cpu_quota | default('') %}
{% set _nice = _profile.nice if _profile.nice is defined else restic_profile_systemd_nice | default('') %}
{% set _io_class = _profile.io_scheduling_class if _profile.io_scheduling_class is defined else restic_profile_systemd_io_scheduling_class | default('') %}
{% set _io_priority = _profile.io_scheduling_priority if _profile.io_scheduling_priority is defined else restic_profile_systemd_io_scheduling_priority | default('') %}
{% if _profile.backup is defined and _profile.retention is defined %}
{%   set _workflow_type = 'backup+retention' %}
{% elif _profile.backup is defined %}
{%   set _workflow_type = 'backup' %}
{% else %}
{%   set _workflow_type = 'retention' %}
{% endif %}
# {{ _unit_name }}.service
# Generated by Ansible - do not edit manually

[Unit]
Description=Restic {{ _workflow_type }} - {{ _profile_name }}
# Soft dependency: start after network is up, but don't require it
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
{% if _cpu_quota != '' %}
CPUQuota={{ _cpu_quota }}
{% endif %}
{% if _nice != '' %}
Nice={{ _nice }}
{% endif %}
{% if _io_class != '' %}
IOSchedulingClass={{ _io_class }}
{% if _io_priority != '' %}
IOSchedulingPriority={{ _io_priority }}
{% endif %}
{% endif %}
ExecStart={{ restic_profile_bin }} {{ _profile_name }} \
    --config {{ restic_profile_config_file }}
# Let systemd capture stdout/stderr in the journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ _unit_name }}

[Install]
WantedBy={{ _restic_profile_install_target | default('multi-user.target') }}
