# roles/restic_profile/templates/etc/systemd/system/restic-profile.timer.j2
{# jinja2: lstrip_blocks: True, trim_blocks: True #}
{% set _profile_name = _profile_item.key %}
{% set _profile = _profile_item.value %}
{% set _is_backup = (_profile.sources | default([]) | length) > 0 %}
{% set _unit_prefix = restic_profile_backup_unit_prefix if _is_backup else restic_profile_forget_unit_prefix %}
{% set _cmd_desc = 'backup' if _is_backup else 'forget' %}
# /etc/systemd/system/{{ _unit_prefix }}{{ _profile_name }}.timer
# Generated by Ansible — do not edit manually

[Unit]
Description=Restic {{ _cmd_desc }} timer — {{ _profile_name }}

[Timer]
# Run on this schedule; Persistent=true catches up a missed run after reboot/sleep
OnCalendar={{ _profile.on_calendar | default('hourly') }}
Persistent=true
# Spread load: fire at a random offset within this window
RandomizedDelaySec={{ _profile.randomized_delay_sec | default('10min') }}

[Install]
WantedBy=timers.target
