# roles/restic_profile/templates/etc/systemd/system/restic-profile.timer.j2
{# jinja2: lstrip_blocks: True, trim_blocks: True #}
{% set _profile_name = item.profile_name %}
{% set _unit_name = item.unit_name %}
{% set _profile = item.profile %}
# {{ _unit_name }}.timer
# Generated by Ansible — do not edit manually

[Unit]
Description=Restic profile timer — {{ _profile_name }}

[Timer]
# Run on this schedule.
# When Persistent=true, a missed run (e.g. system was powered off) is
# triggered immediately after the next boot — keep it false for
# workload-sensitive profiles (e.g. gitlab-backups during business hours).
OnCalendar={{ _profile.on_calendar | default('hourly') }}
Persistent={{ _profile.persistent | default(true) | ansible.builtin.bool | lower }}
# Spread load: fire at a random offset within this window
RandomizedDelaySec={{ _profile.randomized_delay_sec | default('10min') }}

[Install]
WantedBy=timers.target
